# docker/etc/s6-overlay/s6-rc.d/api/run
#!/command/with-contenv sh
set -eu

cd /srv/apps/api
exec bun run serve
