# docker/etc/s6-overlay/s6-rc.d/redis/run
#!/command/with-contenv sh
set -eu

exec redis-server /etc/redis/redis.conf
